# deploy/systemd/omnijob-dedupe.service
[Unit]
Description=OmniJob cross-source dedupe - chained after crawler
Documentation=https://github.com/OM200401/OmniJOB
After=omnijob-crawler.service

[Service]
Type=oneshot
User=om
WorkingDirectory=/home/om/omnijob/apps/api
ExecStart=/home/om/.bun/bin/bun run scripts/dedupe.ts
Environment=QDRANT_URL=http://localhost:6333
StandardOutput=append:/var/log/omnijob-dedupe.log
StandardError=append:/var/log/omnijob-dedupe.log

[Install]
WantedBy=multi-user.target
